# Service_WaterManage/deploy/waterms.service
[Unit]
Description=Water Management System API Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/jhw-ai.com/backend
Environment="PATH=/var/www/jhw-ai.com/backend/.venv/bin"
ExecStart=/var/www/jhw-ai.com/backend/.venv/bin/python -m uvicorn main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
StandardOutput=append:/var/log/waterms/waterms.log
StandardError=append:/var/log/waterms/waterms-error.log

[Install]
WantedBy=multi-user.target